# systemctl/everyday.service
[Unit]
Description="Service called everyday"

[Service]
ExecStart=/usr/bin/everyday.sh
EnvironmentFile=/etc/scripts.env
