# systemctl/gptbible.timer
[Unit]
Description=Run everyday

[Timer]
Unit=gptbible.service
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
