# systemctl/gptbible.timer
[Unit]
Description=Run everyday

[Timer]
OnBootSec=15m
OnUnitActiveSec=24h
Unit=gptbible.service
Persistent=true

[Install]
WantedBy=timers.target
